# .config/systemd/user/actual-gateway.service
[Unit]
Description=Actual Server Gateway (MCP)
After=default.target

[Service]
Type=simple
ExecStart=/usr/bin/java \
    -Xmx2g -XX:+UseZGC -XX:+ExitOnOutOfMemoryError \
    -javaagent:/home/bw/actual/lib/opentelemetry-javaagent.jar \
    -jar /home/bw/actual/build/libs/actual-1.0.0.jar
WorkingDirectory=/home/bw/actual
Environment=DISPLAY=:0
Environment=OTEL_SERVICE_NAME=actual-gateway
Environment=OTEL_EXPORTER_OTLP_ENDPOINT=http://127.0.0.1:4318
Environment=OTEL_EXPORTER_OTLP_PROTOCOL=http/protobuf
Environment=OTEL_METRICS_EXPORTER=otlp
Environment=OTEL_TRACES_EXPORTER=otlp
Environment=OTEL_LOGS_EXPORTER=otlp
Environment=OTEL_INSTRUMENTATION_GRPC_ENABLED=false
Restart=on-failure
RestartSec=5

# Resource limits
MemoryMax=4G
MemoryHigh=3G
CPUQuota=80%
TasksMax=512

# Logging
StandardOutput=null
StandardError=append:/home/bw/vault/data/logs/gateway.log

[Install]
WantedBy=default.target
