# .config/systemd/user/chatterbox.service
[Unit]
Description=Chatterbox TTS Server
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/bw/vault/programs/chatterbox
ExecStart=/home/bw/vault/programs/chatterbox/.venv/bin/python web.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
