# .config/systemd/user/ci-repair.service
[Unit]
Description=Autonomous CI repair — detect and fix GitHub Actions failures
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/bw/actual
ExecStart=/home/bw/actual/scripts/ci-repair.sh
TimeoutStartSec=600
# Inherit user environment for gh auth, PATH (claude, gh, jq)
# Do NOT set GITHUB_TOKEN — script uses gh auth login credentials
Environment=HOME=/home/bw
Environment=PATH=/home/bw/.local/bin:/usr/local/bin:/usr/bin:/bin

[Install]
WantedBy=default.target
