# .config/systemd/user/concierge.service
[Unit]
Description=Claude Concierge Daemon
Documentation=file:///home/bw/.claude/orchestrator/concierge-daemon.py
After=default.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /home/bw/.claude/orchestrator/concierge-daemon.py
WorkingDirectory=/home/bw/.claude/orchestrator
Restart=on-failure
RestartSec=10

# Logging
StandardOutput=append:/home/bw/vault/data/logs/concierge.log
StandardError=append:/home/bw/vault/data/logs/concierge.log

# Resource limits (lightweight daemon)
MemoryMax=256M
CPUQuota=10%

[Install]
WantedBy=default.target
