# .config/systemd/user/dictate-hold.service
[Unit]
Description=Hold-to-dictate (Pause key via evdev)
After=dictation-daemon.service
Requires=dictation-daemon.service

[Service]
Type=simple
Environment=LD_LIBRARY_PATH=%h/.local/share/dictation-venv/lib/python3.14/site-packages/nvidia/cublas/lib
Environment=DISPLAY=:0
ExecStart=%h/bin/dictate-hold
ExecStopPost=/usr/bin/polybar-msg action dictation hook 0
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
