# .config/systemd/user/faster-whisper.service
[Unit]
Description=Faster Whisper transcription daemon
After=pipewire.service

[Service]
Type=simple
Environment=LD_LIBRARY_PATH=%h/.local/share/dictation-venv/lib/python3.14/site-packages/nvidia/cublas/lib
ExecStart=%h/bin/dictation-daemon
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
