# .config/systemd/user/file-version-daemon-v2.service
[Unit]
Description=File Version Daemon v2 (fanotify-based)
Documentation=man:fanotify(7)
After=default.target

[Service]
Type=simple
ExecStart=/home/bw/bin/file-version-daemon-v2
Restart=on-failure
RestartSec=5s

# Resource limits
CPUQuota=20%
MemoryMax=512M
IOWeight=10

# Nice level (low priority)
Nice=10

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/bw/.local/share/file-versions

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=file-version-daemon-v2

[Install]
WantedBy=default.target
