# .config/systemd/user/file-version-daemon.service
[Unit]
Description=File Versioning Daemon
Documentation=file:///home/bw/bin/file-version-daemon
After=default.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /home/bw/bin/file-version-daemon
Restart=on-failure
RestartSec=10

# Logging
StandardOutput=append:/home/bw/vault/data/logs/file-version-daemon.log
StandardError=append:/home/bw/vault/data/logs/file-version-daemon.log

# Resource limits (lightweight daemon)
MemoryMax=512M
CPUQuota=20%

# Nice level (low priority)
Nice=10

[Install]
WantedBy=default.target
