# .config/systemd/user/github-runner.service
[Unit]
Description=GitHub Actions Self-Hosted Runner
After=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/bw/github-actions-runner
ExecStart=/home/bw/github-actions-runner/run.sh
Restart=on-failure
RestartSec=5

# Liberal resource limits for Threadripper PRO
MemoryMax=32G
MemoryHigh=24G
CPUQuota=800%

[Install]
WantedBy=default.target
