# .config/systemd/user/i3-state-writer.service
[Unit]
Description=i3 window state writer for Claude hooks
After=graphical-session.target

[Service]
Type=simple
ExecStart=/home/bw/.claude/hooks/i3_state_writer
Restart=on-failure
RestartSec=5
Environment=DISPLAY=:0

[Install]
WantedBy=default.target
