# .config/systemd/user/kinetic.service
[Unit]
Description=Kinetic MCP Gateway
Documentation=file:///home/bw/vault/docs/specs/kinetic-vocabulary.md
After=default.target

[Service]
Type=simple
ExecStart=/usr/bin/java -Xmx2g -XX:+UseZGC -XX:+ExitOnOutOfMemoryError -jar /home/bw/vault/programs/kinetic/build/libs/kinetic-1.0.0.jar
WorkingDirectory=/home/bw/vault/programs/kinetic
Environment=DISPLAY=:0
Restart=on-failure
RestartSec=5

# Resource limits
MemoryMax=3G
MemoryHigh=2G
CPUQuota=80%
TasksMax=200

# Logging (stderr only — stdout is MCP JSON-RPC)
StandardOutput=null
StandardError=append:/home/bw/vault/data/logs/kinetic.log

[Install]
WantedBy=default.target
