# .config/systemd/user/openclaw.service
[Unit]
Description=OpenClaw Gateway
After=network.target

[Service]
Type=simple
ExecStart=/home/bw/.local/bin/openclaw gateway
WorkingDirectory=/home/bw
Environment=OLLAMA_API_KEY=ollama-local
Environment=PATH=/home/bw/.local/bin:/usr/local/bin:/usr/bin:/bin
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
