# .config/systemd/user/overnight-review.service
[Unit]
Description=Overnight Frontier Validation Review (idle-triggered)
After=actual-gateway.service

[Service]
Type=oneshot
Environment=HOME=/home/bw
ExecStart=/home/bw/.local/bin/overnight-review-check.sh
WorkingDirectory=/home/bw/actual
TimeoutStopSec=7200
Nice=19
IOSchedulingClass=idle
