# .config/systemd/user/overnight-review.timer
[Unit]
Description=Check idle time and run overnight review if idle 6+ hours

[Timer]
OnBootSec=30min
OnUnitActiveSec=15min
Persistent=false

[Install]
WantedBy=timers.target
