# .config/systemd/user/time-machine-archive.service
[Unit]
Description=Time Machine Archive (daily restic backup)
Documentation=file:///home/bw/vault/programs/time-machine/docs/PRD.md
After=network.target time-machine-snapshot.service

[Service]
Type=oneshot
ExecStart=/usr/bin/python3 /home/bw/vault/programs/time-machine/src/archive-manager.py
StandardOutput=journal
StandardError=journal

# Resource limits (lower priority than snapshots)
CPUQuota=50%
Nice=19
IOWeight=50

# Timeout (archives can take hours)
TimeoutStartSec=3h

[Install]
WantedBy=multi-user.target
