# .config/systemd/user/time-machine-snapshot.timer
[Unit]
Description=Time Machine Snapshot Timer (every 1 minute)
Documentation=file:///home/bw/vault/programs/time-machine/docs/PRD.md

[Timer]
OnCalendar=*-*-* *:*:00
OnBootSec=2m
Persistent=true
AccuracySec=1s

[Install]
WantedBy=timers.target
