# .config/systemd/user/time-machine-archive.timer
[Unit]
Description=Time Machine Archive Timer (daily at 3am)
Documentation=file:///home/bw/vault/programs/time-machine/docs/PRD.md

[Timer]
OnCalendar=*-*-* 03:00:00
Persistent=true
RandomizedDelaySec=600

[Install]
WantedBy=timers.target
