# .config/systemd/user/time-machine-snapshot.service
[Unit]
Description=Time Machine Snapshot
Documentation=file:///home/bw/vault/programs/time-machine/docs/PRD.md
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/bin/python3 /home/bw/vault/programs/time-machine/src/snapshot-manager.py
StandardOutput=journal
StandardError=journal

# Resource limits (low priority background task)
CPUQuota=10%
Nice=15
IOWeight=100

# Timeout (snapshots should complete in <5 minutes)
TimeoutStartSec=300

[Install]
WantedBy=multi-user.target
