# .config/systemd/user/vault-rag-watcher.service
[Unit]
Description=Vault RAG File Watcher
Documentation=https://github.com/your/vault-rag
After=default.target

[Service]
Type=simple
ExecStart=/home/bw/vault/programs/vault-rag/.venv/bin/python /home/bw/vault/programs/vault-rag/scripts/watcher.py
WorkingDirectory=/home/bw/vault/programs/vault-rag
Environment=VAULT_ROOT=/home/bw/vault
Restart=on-failure
RestartSec=10

# Logging
StandardOutput=append:/home/bw/vault/data/logs/watcher.log
StandardError=append:/home/bw/vault/data/logs/watcher.log

[Install]
WantedBy=default.target
